# root/etc/services.d/ttrss-daemon/run
#!/usr/bin/with-contenv sh

cd /var/www/ttrss
exec s6-setuidgid www-data php -f /var/www/ttrss/update.php -- --update
