# root/etc/services.d/ttrss-daemon/run
#!/usr/bin/with-contenv sh

while true; do
  cd /var/www/ttrss
  exec s6-setuidgid www-data php -f /var/www/ttrss/update_daemon2.php
  sleep 5m
done
